# dotconfig/containers/systemd/orcaslicer.container
[Unit]
Description=OrcaSlicer 3D printing slicer with web interface
After=network-online.target
Wants=network-online.target

[Container]
Image=quay.io/ckyrouac/orcaslicer:latest
PublishPort=3001:3001
Environment=GDK_SCALE=2
Environment=GDK_DPI_SCALE=0.50
Environment=PUID=1000
Environment=PGID=1000
# NVIDIA environment variables to reduce GPU context contention
Environment=__GL_YIELD=USLEEP
Environment=__GL_SYNC_TO_VBLANK=0
Environment=__GL_THREADED_OPTIMIZATIONS=0
Environment=__GL_GSYNC_ALLOWED=0
Environment=__GL_VRR_ALLOWED=0
Volume=%h/.config/orcaslicer:/config
Volume=%h/projects/3dprinter:/3dprinter:Z
# Override Xvfb startup to disable indirect GLX (reduces NVIDIA context contention)
Volume=%h/.config/container-overrides/svc-xorg-run:/etc/s6-overlay/s6-rc.d/svc-xorg/run:ro
PodmanArgs=--shm-size=1gb --gpus all --env-file %t/orcaslicer.env
AutoUpdate=registry

[Service]
Restart=always
TimeoutStartSec=900
LoadCredential=cf_token:%h/.secrets/orcaslicer-cf-token
ExecStartPre=/bin/sh -c 'umask 077 && printf "CF_TOKEN=" > %t/orcaslicer.env && cat %d/cf_token >> %t/orcaslicer.env'
ExecStopPost=/bin/rm -f %t/orcaslicer.env

[Install]
WantedBy=default.target
